# deploy/ehc-vllm.service
[Unit]
Description=EHC vLLM Shared Inference Service
After=network.target

[Service]
Environment="VLLM_PLUGINS="
Type=simple
User=phungkien
WorkingDirectory=/home/phungkien/EHC_HELPDESK
Environment="PATH=/home/phungkien/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/usr/bin/python3 -m vllm.entrypoints.openai.api_server \
    --model Qwen/Qwen3-8B \
    --host 0.0.0.0 \
    --port 8000 \
    --enable-reasoning \
    --reasoning-parser deepseek_r1 \
    --dtype half \
    --tensor-parallel-size 2 \
    --max-model-len 32768 \
    --gpu-memory-utilization 0.85
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
